# launchd/com.sectiontracker.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sectiontracker.update</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/wmeleis/committees/nu-docs/Curriculum/SectionTracker/update.sh</string>
  </array>
  <!-- Fire every 30 min (+ on load/wake). update.sh gates itself to one
       SUCCESSFUL refresh per day starting ~6am, retrying after any failure
       (e.g. network not up at 6:30am) instead of leaving the site stale. -->
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/wmeleis/committees/nu-docs/Curriculum/SectionTracker/data/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/wmeleis/committees/nu-docs/Curriculum/SectionTracker/data/launchd.log</string>
</dict>
</plist>
